# apcupsd2mqtt-php/rootfs/etc/s6-overlay/s6-rc.d/apcupsd2mqtt-php-srv/run
#!/bin/sh
exec php /opt/apcupsd2mqtt-php/start.php